# compact/supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod --config /etc/mongodb.conf

[program:webgme]
directory=/webgme
command=/usr/bin/npm start